# backend/deploy/ecs/viralforge-backend.service
[Unit]
Description=ViralForge Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=viralforge
Group=viralforge
WorkingDirectory=/opt/viralforge/backend
EnvironmentFile=/etc/viralforge/backend.env
ExecStart=/usr/bin/node /opt/viralforge/backend/src/main.mjs
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/viralforge

[Install]
WantedBy=multi-user.target
